# setup/systemd/skillops-backup.service
[Unit]
Description=SkillOps Backup
Documentation=https://github.com/your-org/SkillOps/blob/main/docs/OPERATIONS.md
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h
ExecStart=/usr/bin/env bash %h/SkillOps/setup/backup/backup.sh
EnvironmentFile=-%h/.config/skillops/skillops.env

[Install]
WantedBy=default.target
