# setup/systemd/skillops-backup.timer
[Unit]
Description=SkillOps Backup Timer
Documentation=https://github.com/your-org/SkillOps/blob/main/docs/OPERATIONS.md

[Timer]
OnCalendar=*-*-* 01:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
